# configs/zabbix_server/supervisord.conf
[supervisord]
nodaemon=true

[program:zabbix-server]
command=/usr/sbin/zabbix_server -c /etc/zabbix/zabbix_server.conf
autostart=true
autorestart=true
user=zabbix
